# profile.d/node.sh
if [ -d /opt/node ]; then
  PATH="/opt/node/bin:$PATH"
fi
