# profile.d/rcm.sh
PATH="$HOME/.dotfiles/rcm/bin:$PATH"
MANPATH="$HOME/.dotfiles/rcm/share/man:$MANPATH"
