# profile.d/stty.sh
[ -z "$PS1" ] && return # interactive only

# disable tty stop and start
stty stop ''
stty start ''
